# provd/debian/postinst
#!/bin/sh

set -e

systemctl enable provd.service
systemctl start provd.service

exit 0
