# provd/systemd/provd.service
[Unit]
Description=Ubuntu Desktop Provisioning Daemon
Requires=provd.socket

[Service]
ExecStart=/usr/bin/provd
Restart=on-failure
ProtectSystem=strict
PrivateTmp=true
